# silabs-multiprotocol/rootfs/etc/s6-overlay/s6-rc.d/hassio-discovery/finish
#!/usr/bin/with-contenv bashio
#==============================================================================
# hassio-discoveryfinish script
#==============================================================================
bashio::log.info "hassio-discovery ended with exit code ${1} (signal ${2})..."
