# silabs-multiprotocol/rootfs/etc/s6-overlay/s6-rc.d/zeroconf/finish
#!/usr/bin/with-contenv bashio
#==============================================================================
# Zeroconf broadcaster finish script
#==============================================================================
bashio::log.info "zeroconf ended with exit code ${1} (signal ${2})..."
